# app/dos.cron
0 */3 * * * /etc/weather.rc
#0 * * * * /etc/currency.rc
